# conf/cloud9.conf

[program:cloud9]
command = /home/user/.c9/node/bin/node /home/user/c9sdk/server.js --listen 0.0.0.0 --port 8080 -w /workspace
autostart = true
autorestart = true
environment = NODE_ENV="production"
